# docker_logger/logger.crontab
* * * * * bash /app/docker_logger/docker_logs_to_matrix.sh
