# docker_logger/logger.crontab
*/5 * * * * bash /app/docker_logger/scripts/docker_logs_to_matrix.sh > /proc/1/fd/1 2>/proc/1/fd/2
