# scripts/log_directory_to_matrix.crontab
# assuming log.py from this project is in cwd
PASSWORD='abc'
* * * * * bash /some/path/log_directory_to_matrix.sh
