# stats/logger.crontab
0 * * * * bash /app/stats/scripts/log_stats.sh > /proc/1/fd/1 2>/proc/1/fd/2
